# .docker/config/supervisor/supervisor.ini
[supervisord]
serverurl=unix:///var/run/supervisor.sock
user=root
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=900
autorestart=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
